# com.healthcare.intelligence.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.healthcare.intelligence</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>REPLACE_WITH_FULL_PATH/start_permanent_service.sh</string>
    </array>
    
    <key>WorkingDirectory</key>
    <string>REPLACE_WITH_FULL_PATH</string>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>
    
    <key>StandardOutPath</key>
    <string>REPLACE_WITH_FULL_PATH/healthcare_service.log</string>
    
    <key>StandardErrorPath</key>
    <string>REPLACE_WITH_FULL_PATH/healthcare_service_error.log</string>
    
    <key>ThrottleInterval</key>
    <integer>10</integer>
    
</dict>
</plist> 